# runit/mullvad-daemon/finish
#!/bin/sh
set -e

exec /lib/runit/finish-exec "$@"
